# hdf5-datalogger/rootfs/etc/services.d/hdf5_datalogger/run
#!/usr/bin/with-contenv bash
set -e

# parent del package, NON la cartella del package
export PYTHONPATH="/usr/lib:${PYTHONPATH}"

read_interval() {
  if [ -f /data/options.json ]; then
    local val
    val=$(jq -r '.update_interval // 60' /data/options.json 2>/dev/null || echo 60)
    echo "$val" | grep -Eq '^[0-9]+$' && echo "$val" || echo 60
  else
    echo 60
  fi
}

log_startup_summary() {
  local now ts output interval domains_raw

  now=$(date -u +"%Y-%m-%dT%H:%M:%SZ")

  if [ -f /data/options.json ]; then
    output=$(jq -r '.output_path // "/share/example_addon_output.txt"' /data/options.json 2>/dev/null || echo "/share/example_addon_output.txt")
    interval=$(jq -r '.update_interval // 60' /data/options.json 2>/dev/null || echo 60)

    domains_raw=$(jq -r '.include_domains // [] | join(", ")' /data/options.json 2>/dev/null || echo "")
  else
    output="/share/example_addon_output.txt"
    interval="60"
    domains_raw=""
  fi

  echo "[INFO] ==============================="
  echo "[INFO] HDF5 DataLogger avviato"
  echo "[INFO] Timestamp avvio (UTC): ${now}"
  echo "[INFO] Output file configurato: ${output}"
  echo "[INFO] Update interval: ${interval} secondi"
  if [ -n "${domains_raw}" ]; then
    echo "[INFO] include_domains (config): ${domains_raw}"
  else
    echo "[INFO] include_domains (config): (vuoto => domini di default)"
  fi
  echo "[INFO] ==============================="
}

log_startup_summary

while true; do
  python3 /usr/bin/logger.py || true
  sleep "$(read_interval)"
done
