# hdf5-datalogger/rootfs/etc/services.d/hdf5_datalogger/run
#!/usr/bin/with-contenv bash
set -e

# parent del package, NON la cartella del package
export PYTHONPATH="/usr/lib:${PYTHONPATH}"

read_interval() {
  if [ -f /data/options.json ]; then
    local val
    val=$(jq -r '.update_interval // 60' /data/options.json 2>/dev/null || echo 60)
    echo "$val" | grep -Eq '^[0-9]+$' && echo "$val" || echo 60
  else
    echo 60
  fi
}

while true; do
  python3 /usr/bin/logger.py || true
  sleep "$(read_interval)"
done
